# dot_Library/LaunchAgents/local.mai-model-serve.plist.tmpl
{{- if and (eq .chezmoi.os "darwin") (has .chezmoi.hostname .hosts.ai_machines) -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.mai-model-serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/mai-model</string>
        <string>serve</string>
        <string>--bind</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>7860</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash, but back off so a broken plist doesn't spin -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <!-- mai-model resolves $HOME for paths; launchd already sets HOME for
         a user-domain agent. PATH is set so the bash dispatcher can find
         basic utilities (date, find, etc.) without depending on a shell
         rc file. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/mai-model-serve.out.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/mai-model-serve.err.log</string>
</dict>
</plist>
{{ end -}}
